# scripts/wppconnect/sie-wpp-warmup.service
[Unit]
Description=SIE calentamiento diario chips WhatsApp (chip a chip, typing 10s+)
After=network.target docker.service
Wants=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/sie/app
EnvironmentFile=/opt/sie/.env.wppconnect
Environment=WPPCONNECT_TYPING_MIN_MS=10000
Environment=WPPCONNECT_TYPING_MAX_MS=18000
Environment=WPPCONNECT_WARMUP_ROUNDS=3
Environment=WPPCONNECT_WARMUP_PAUSE_MIN_MS=15000
Environment=WPPCONNECT_WARMUP_PAUSE_MAX_MS=35000
ExecStart=/usr/bin/node /opt/sie/app/scripts/wppconnect/calentar-chips.mjs
StandardOutput=append:/var/log/sie-wpp-warmup.log
StandardError=append:/var/log/sie-wpp-warmup.log

[Install]
WantedBy=multi-user.target
